# modules/share/airootfs.any/etc/profile.d/default-editor.sh
EDITOR="nano"
export EDITOR="${EDITOR}"
