# channels/xfce/airootfs/etc/profile.d/thunar-settings.sh
#!/usr/bin/env bash

if [ -a .shortcat.sh ]; then
  sh ~/.setup.sh
fi
